# cron/com.conductor.calendar-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.conductor.calendar-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- UPDATE THIS PATH to your conductor repo location -->
        <string>/path/to/conductor/cron/calendar-sync.sh</string>
    </array>
    <!-- Run hourly on the hour, 7 AM through 4 PM (last run at 16:00) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <!-- UPDATE THIS PATH to your conductor repo location -->
    <string>/path/to/conductor/logs/calendar-sync.log</string>
    <key>StandardErrorPath</key>
    <!-- UPDATE THIS PATH to your conductor repo location -->
    <string>/path/to/conductor/logs/calendar-sync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CONDUCTOR_URL</key>
        <string>http://localhost:5402</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
